# opentelemetry-collector/rootfs/etc/s6-overlay/s6-rc.d/opentelemetry-collector-log/run
#!/usr/bin/env bash
exec logutil-prefix opentelemetry-collector
