# cadvisor/rootfs/etc/s6-overlay/s6-rc.d/cadvisor/run
#!/command/execlineb -P
cadvisor \
	--v=1 \
	--docker_only \
	--housekeeping_interval=15s \
	--update_machine_info_interval=15s
