# cadvisor/rootfs/etc/s6-overlay/s6-rc.d/cadvisor/run
#!/usr/bin/env bash
source /var/cadvisor/default
exec 2>&1
exec cadvisor \
		--docker_only \
		--housekeeping_interval=30s
