# promtail/rootfs/etc/s6-overlay/s6-rc.d/promtail/run
#!/usr/bin/env bash
exec 2>&1
exec promtail \
		--config.expand-env=true \
		--config.file=/etc/promtail/config.yml \
		--server.enable-runtime-reload
