# grafana/rootfs/etc/s6-overlay/s6-rc.d/grafana/run
# This will be replaced during the build process
# COPY --from=grafana /run.sh /etc/s6-overlay/s6-rc.d/grafana/run
# See Dockerfile
