# grafana/rootfs/etc/s6-overlay/s6-rc.d/grafana/run
#!/usr/bin/env bash
# shellcheck disable=SC1091
source /var/grafana/default
grafana-run
