# node-exporter/rootfs/etc/s6-overlay/s6-rc.d/node-exporter/run
#!/command/execlineb -P
node_exporter \
	--path.rootfs=/rootfs/root \
	--path.sysfs=/rootfs/sys \
	--path.procfs=/rootfs/proc \
	--path.udev.data=/rootfs/root/run/udev/data \
	--collector.systemd \
	--collector.processes \
	--collector.textfile.directory=/etc/node-exporter \
	--collector.filesystem.mount-points-exclude=^/\(sys\|proc\|dev\|host\|etc\)\(\$\|/\) \
	--no-collector.ipvs
