# alloy/rootfs/etc/s6-overlay/s6-rc.d/alloy-log/run
#!/usr/bin/env bash
exec logutil-prefix alloy
