# pyroscope/rootfs/etc/s6-overlay/s6-rc.d/pyroscope/run
#!/command/execlineb -P
pyroscope \
	--config.file=/etc/pyroscope/config.yaml
