# prometheus/rootfs/etc/s6-overlay/s6-rc.d/prometheus-log/run
#!/usr/bin/env bash
exec logutil-prefix prometheus
