# grafana/rootfs/etc/s6-overlay/s6-rc.d/grafana-log/run
#!/usr/bin/env bash
exec logutil-prefix grafana
