# prometheus/rootfs/etc/s6-overlay/s6-rc.d/prometheus/run
#!/command/execlineb -P
prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/data/prometheus
