# prometheus/rootfs/etc/s6-overlay/s6-rc.d/prometheus/run
#!/command/execlineb -P
prometheus \
	--config.file=/etc/prometheus/prometheus.yml
	--enable-feature=exemplar-storage \
	--enable-feature=native-histograms \
	--storage.tsdb.path=/data/prometheus \
	--web.enable-lifecycle \
	--web.enable-otlp-receiver \
	--web.enable-remote-write-receiver \
	--log.level=info
