# opentelemetry-collector/rootfs/etc/s6-overlay/s6-rc.d/opentelemetry-collector/run
#!/command/execlineb -P
otelcol-contrib \
	--feature-gates service.profilesSupport \
	--config=file:/etc/otelcol/config.yaml
