# loki/rootfs/etc/s6-overlay/s6-rc.d/loki/run
#!/command/execlineb -P
loki \
	--config.file=/etc/loki/config.yaml
