# blackbox-exporter/rootfs/etc/s6-overlay/s6-rc.d/blackbox-exporter/run
#!/command/execlineb -P
blackbox_exporter --config.file=/etc/blackbox_exporter/config.yml
