# tempo/rootfs/etc/s6-overlay/s6-rc.d/tempo/run
#!/command/execlineb -P
tempo \
	--config.file=/etc/tempo/config.yaml
