# tempo/rootfs/etc/s6-overlay/s6-rc.d/tempo/run
#!/usr/bin/env bash
exec 2>&1
exec tempo \
	--config.file=/etc/tempo/config.yaml
